# autostart/lxqt-autostart.desktop
[Desktop Entry]
Exec=autostart
Name=autostart
OnlyShowIn=LXQt;
Type=Application
Version=1.0
