# autostart/lxqt-autostart.desktop
[Desktop Entry]
Exec=lxqt-autostart
Name=lxqt-autostart
OnlyShowIn=LXQt;
Type=Application
Version=1.0
