# services/sway-graphical-session.target
[Unit]
Description=User Graphical Session
Requires=graphical-session.target
After=graphical-session.target
#AllowIsolate=yes
